# deploy/storybot-kiosk.service
# DEPRECATED: This systemd service is superseded by the GNOME autostart file
# at /home/ari/.config/autostart/storybot-kiosk.desktop
# The autostart approach handles DISPLAY and XAUTHORITY automatically.
# This file is kept for reference only and is NOT installed by install.sh.
[Unit]
Description=StoryBot Chromium Kiosk
After=graphical.target storybot.service
Wants=storybot.service

[Service]
Type=simple
User=storybot
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/storybot/.Xauthority"
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/chromium-browser --noerrdialogs --disable-infobars --disable-session-crashed-bubble --disable-translate --no-first-run --start-fullscreen --kiosk http://localhost:8000
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
