# deploy/storybot.service
# Template: __INSTALL_USER__ and __INSTALL_DIR__ are substituted by
# deploy/install.sh (from INSTALL_USER in .env) when this unit is installed.
[Unit]
Description=StoryBot FastAPI Server
After=network-online.target systemd-udev-settle.service storybot-nfc-reset.service pcscd.service llama-server.service
Wants=network-online.target systemd-udev-settle.service storybot-nfc-reset.service llama-server.service
Requires=pcscd.service

[Service]
Type=simple
User=__INSTALL_USER__
Group=__INSTALL_USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=-__INSTALL_DIR__/.env
Environment="PATH=__INSTALL_DIR__/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/usr/bin/udevadm settle
ExecStartPre=/bin/sleep 4
ExecStartPre=/usr/local/bin/storybot-rollback.sh
ExecStart=__INSTALL_DIR__/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
