# deploy/storybot.service
[Unit]
Description=StoryBot FastAPI Server
After=network.target

[Service]
Type=simple
User=ari
Group=ari
WorkingDirectory=/home/ari/storybot
Environment="PATH=/home/ari/storybot/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ari/storybot/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
